# git-autocommit.service
[Unit]
Description=git autocommit

[Service]
ExecStart=/home/serge/.git-autocommit

[Install]
WantedBy=default.target
